# deployments/systemd/gomcp-server.service
[Unit]
Description=GoMCP Multi-Tool Example Server
Documentation=https://github.com/localrivet/gomcp
After=network.target

[Service]
# --- Configuration ---
# IMPORTANT: Replace '/opt/gomcp-server' with the actual deployment directory on your server.
WorkingDirectory=/opt/gomcp-server

# Build the executable and place it in the WorkingDirectory.
# Example build command (run from repo root):
# CGO_ENABLED=0 GOOS=linux go build -o ./deployments/systemd/multi-tool-server ./examples/basic/server/*.go
# Then copy 'multi-tool-server' and 'start-gomcp-server.sh' to '/opt/gomcp-server' on the target machine.

# Executable path (relative to WorkingDirectory)
ExecStart=/opt/gomcp-server/start-gomcp-server.sh

# --- Security ---
# Recommended: Create a dedicated user/group for the service
# useradd --system --shell /sbin/nologin gomcp-user
# User=gomcp-user
# Group=gomcp-user
# For simplicity in example, using 'nobody'. Change this for production.
User=nobody
Group=nogroup

# --- Process Management ---
Type=simple
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

# --- Logging ---
# Redirect stdout and stderr to systemd-journald
StandardOutput=journal
StandardError=journal

# --- Environment (Optional) ---
# If the server needed environment variables (like an API key),
# you could set them here or preferably in the start script.
# Environment="MCP_API_KEY=your_key_here"
# EnvironmentFile=/etc/gomcp-server/environment

[Install]
WantedBy=multi-user.target